# examples/dashboard/deploy/dashboard.service
[Unit]
Description=Dashboard HTTP service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/dashboard/dashboard --serve --bind 127.0.0.1 --port 6743 --tier flash
WorkingDirectory=/opt/dashboard
EnvironmentFile=/etc/ai-ledger.env
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=yes
ReadWritePaths=/var/lib/ai-ledger
User=root
Group=root

[Install]
WantedBy=multi-user.target
